# supervisor.conf
[program:myproject.com.br]
command=/home/vagrant/www/Django/start.sh
user=vagrant
stdout_logfile=/home/vagrant/www/Django/logs/supervisor.log
redirect_stderr=true
autostart=true
autorestart=true
